# micromonitor/micromonitor.service
[Unit]
Description=MicroMonitor System Monitoring Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/claude-experiments/micromonitor
ExecStart=/usr/bin/node /root/claude-experiments/micromonitor/src/server.js
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
Environment=NODE_ENV=production
Environment=PORT=3000

[Install]
WantedBy=multi-user.target